# homebridge/supervisord.conf
[supervisord]
nodaemon=true

[program:homebridge]
command=/etc/init.d/homebridge.sh start
