# homebridge/supervisord.conf
[supervisord]
nodaemon=true

[program:homebridge]
command=homebridge start
